# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mongodb]
command=mongod --bind_ip 127.0.0.1 --dbpath /var/lib/mongodb
user=root
autostart=true
autorestart=true
stdout_logfile=/var/log/mongodb.log
stderr_logfile=/var/log/mongodb.err.log

[program:auth-service]
command=node /app/backend/auth-service/server.js
directory=/app/backend/auth-service
environment=PORT="5001"
autostart=true
autorestart=true
stdout_logfile=/var/log/auth-service.log
stderr_logfile=/var/log/auth-service.err.log

[program:labour-service]
command=node /app/backend/labour-service/server.js
directory=/app/backend/labour-service
environment=PORT="5002"
autostart=true
autorestart=true
stdout_logfile=/var/log/labour-service.log
stderr_logfile=/var/log/labour-service.err.log

[program:attendance-service]
command=node /app/backend/attendance-service/server.js
directory=/app/backend/attendance-service
environment=PORT="5003"
autostart=true
autorestart=true
stdout_logfile=/var/log/attendance-service.log
stderr_logfile=/var/log/attendance-service.err.log

[program:deployment-service]
command=node /app/backend/deployment-service/server.js
directory=/app/backend/deployment-service
environment=PORT="5004"
autostart=true
autorestart=true
stdout_logfile=/var/log/deployment-service.log
stderr_logfile=/var/log/deployment-service.err.log

[program:reporting-service]
command=node /app/backend/reporting-service/server.js
directory=/app/backend/reporting-service
environment=PORT="5005"
autostart=true
autorestart=true
stdout_logfile=/var/log/reporting-service.log
stderr_logfile=/var/log/reporting-service.err.log

[program:gateway]
command=node /app/backend/gateway/server.js
directory=/app/backend/gateway
environment=PORT="5000"
autostart=true
autorestart=true
stdout_logfile=/var/log/gateway.log
stderr_logfile=/var/log/gateway.err.log

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
